# com.claude.cli-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude.cli-monitor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>REPLACE_WITH_HOME/.claude/monitors/cli-updates.py</string>
    </array>

    <!-- Run at 8am, 2pm, and 8pm daily -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>14</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>20</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>REPLACE_WITH_HOME/.claude/monitors/cli-updates-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>REPLACE_WITH_HOME/.claude/monitors/cli-updates-stderr.log</string>

    <key>KeepAlive</key>
    <false/>

    <key>WorkingDirectory</key>
    <string>REPLACE_WITH_HOME/.claude/monitors</string>
</dict>
</plist>
